# install_scripts/startup/start-containers.service
[Unit]
Description=start docker containers

[Service]
ExecStart=/usr/local/bin/start-containers.sh start

[Install]
WantedBy=multi-user.target
